# deploy/systemd/promo_v2-daily-backup.timer.example
[Unit]
Description=Run Promo V2 daily backup

[Timer]
OnCalendar=*-*-* 02:15:00
Persistent=true
Unit=promo_v2-daily-backup.service

[Install]
WantedBy=timers.target
